# service/nghttpx/finish
#!/bin/sh
if [[ "$1" -eq 0 ]]
then s6-svscanctl -rt /service
else s6-svscanctl -pt /service
fi
